# test.service
[Unit]
Description=test

[Service]
ExecStart=/usr/local/bin/test.sh
Type=onshot
